# debian/wb-homa-rfsniffer.postrm
#!/bin/bash

CONFFILE="/etc/wb-homa-rfsniffer.conf"

if [ "$1" = "purge" ]; then
    if [ -x /usr/bin/ucf ]; then
        ucf --purge $CONFFILE
    fi
fi

update-rc.d -f wb-homa-rfsniffer remove
