# debian/wb-homa-rfsniffer.postinst
#!/bin/bash

update-rc.d wb-homa-rfsniffer defaults
invoke-rc.d wb-homa-rfsniffer start
